# tradingagents-celery-heavy.service
[Unit]
Description=TradingAgents Celery Heavy Worker (batch / backfill tasks)
After=network.target redis-server.service

[Service]
User=admin
WorkingDirectory=/opt/tradingagents
EnvironmentFile=/opt/tradingagents/.env.prod
Environment=PYTHONUTF8=1
ExecStart=/opt/tradingagents/venv/bin/celery \
    -A server.celery_app worker \
    --loglevel=info \
    --concurrency=1 \
    --queues=heavy
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
